# lib/udev/rules.d/85-brltty.rules
# udev rules file for brltty
#

ACTION!="add", GOTO="brltty_rules_end"
SUBSYSTEM!="usb_device", GOTO="brltty_rules_end"

# Alva
SYSFS{idVendor}=="06b0", SYSFS{idProduct}=="0001", RUN+="/lib/brltty/brltty.sh -b al -d usb:"
SYSFS{idVendor}=="0798", SYSFS{idProduct}=="0640", RUN+="/lib/brltty/brltty.sh -b al -d usb:"
SYSFS{idVendor}=="0798", SYSFS{idProduct}=="0680", RUN+="/lib/brltty/brltty.sh -b al -d usb:"

# Baum
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="fe70", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="fe71", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="fe72", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="fe73", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="fe74", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="fe75", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="fe76", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="fe77", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2000", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2001", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2002", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2007", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2008", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2009", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2010", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2011", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2012", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2014", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2015", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"
SYSFS{idVendor}=="0904", SYSFS{idProduct}=="2016", RUN+="/lib/brltty/brltty.sh -b bm -d usb:"

# BrailleSense
SYSFS{idVendor}=="045e", SYSFS{idProduct}=="930a", RUN+="/lib/brltty/brltty.sh -b hm -d usb:"

# EuroBraille
SYSFS{idVendor}=="c251", SYSFS{idProduct}=="1122", RUN+="/lib/brltty/brltty.sh -b eu -d usb:"
SYSFS{idVendor}=="c251", SYSFS{idProduct}=="1124", RUN+="/lib/brltty/brltty.sh -b eu -d usb:"

# FreedomScientific
SYSFS{idVendor}=="0f4e", SYSFS{idProduct}=="0100", RUN+="/lib/brltty/brltty.sh -b fs -d usb:"
SYSFS{idVendor}=="0f4e", SYSFS{idProduct}=="0111", RUN+="/lib/brltty/brltty.sh -b fs -d usb:"
SYSFS{idVendor}=="0f4e", SYSFS{idProduct}=="0112", RUN+="/lib/brltty/brltty.sh -b fs -d usb:"

# HandyTech
SYSFS{idVendor}=="0921", SYSFS{idProduct}=="1200", RUN+="/lib/brltty/brltty.sh -b ht -d usb:"
#SYSFS{idVendor}=="0403", SYSFS{idProduct}=="6001", RUN+="/lib/brltty/brltty.sh -b ht -d usb:"

# Papenmeier
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="f208", RUN+="/lib/brltty/brltty.sh -b pm -d usb:"

# Voyager
SYSFS{idVendor}=="0798", SYSFS{idProduct}=="0001", RUN+="/lib/brltty/brltty.sh -b vo -d usb:"


LABEL="brltty_rules_end"
