# var/lib/dpkg/info/anacron.conffiles
/etc/apm/event.d/anacron
/etc/acpi/resume.d/85-anacron.sh
/etc/acpi/ac.d/85-anacron.sh
/etc/acpi/suspend.d/15-anacron.sh
/etc/acpi/battery.d/15-anacron.sh
/etc/init.d/anacron
/etc/cron.d/anacron
/etc/cron.monthly/0anacron
/etc/cron.weekly/0anacron
/etc/cron.daily/0anacron
/etc/anacrontab
