# var/lib/dpkg/info/apparmor.preinst
#!/bin/sh
# preinst script for apparmor
#
# see: dh_installdeb(1)
set -e

# Remove symlinks in rc[\d].d.

if dpkg --compare-versions "${2}" le-nl "2.0.1+510.dfsg-0ubuntu7"
then
	update-rc.d -f apparmor remove
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
