# var/lib/dpkg/info/apport.conffiles
/etc/logrotate.d/apport
/etc/default/apport
/etc/init.d/apport
/etc/cron.daily/apport
/etc/apport/blacklist.d/apport
