# var/lib/dpkg/info/apport.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/apport
/usr/share/doc/apport/README.Debian
/usr/share/doc/apport/copyright
/usr/share/doc/apport/examples
/usr/share/doc/apport/examples/Makefile
/usr/share/doc/apport/examples/README
/usr/share/doc/apport/examples/libapport.c.gz
/usr/share/doc/apport/crashdb-conf.txt
/usr/share/doc/apport/package-hooks.txt
/usr/share/doc/apport/data-format.pdf.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/apport-cli.1.gz
/usr/share/man/man1/apport-collect.1.gz
/usr/share/man/man1/apport-unpack.1.gz
/usr/share/man/man1/ubuntu-bug.1.gz
/usr/share/apport
/usr/share/apport/apport
/usr/share/apport/apport-checkreports
/usr/share/apport/package_hook
/usr/share/apport/kernel_crashdump
/usr/share/apport/kernel_oops
/usr/share/apport/gcc_ice_hook
/usr/share/apport/apportcheckresume
/usr/share/apport/testsuite
/usr/share/apport/testsuite/test-apport
/usr/share/apport/testsuite/test-hooks
/usr/share/apport/testsuite/run-tests
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_apport.py
/usr/share/apport/package-hooks/source_linux.py
/usr/share/apport/general-hooks
/usr/share/apport/general-hooks/ubuntu.py
/usr/share/apport/general-hooks/automatix.py
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/scalable
/usr/share/icons/hicolor/scalable/apps
/usr/share/icons/hicolor/scalable/apps/apport.svg
/usr/share/icons/hicolor/scalable/mimetypes
/usr/share/mime
/usr/share/mime/packages
/usr/share/mime/packages/apport.xml
/usr/share/pyshared-data
/usr/share/pyshared-data/apport
/usr/bin
/usr/bin/apport-unpack
/usr/bin/apport-cli
/usr/bin/ubuntu-bug
/usr/bin/apport-collect
/usr/lib
/usr/lib/python2.5
/usr/lib/python2.5/site-packages
/usr/lib/python2.6
/usr/lib/python2.6/dist-packages
/etc
/etc/logrotate.d
/etc/logrotate.d/apport
/etc/default
/etc/default/apport
/etc/init.d
/etc/init.d/apport
/etc/cron.daily
/etc/cron.daily/apport
/etc/apport
/etc/apport/blacklist.d
/etc/apport/blacklist.d/apport
/usr/share/doc/apport/TODO
/usr/share/doc/apport/AUTHORS
/usr/share/doc/apport/changelog.Debian.gz
/usr/share/icons/hicolor/scalable/mimetypes/text-x-apport.svg
