# var/lib/dpkg/info/apt-xapian-index.list
/.
/etc
/etc/cron.weekly
/etc/cron.weekly/apt-xapian-index
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/update-apt-xapian-index.8.gz
/usr/share/doc
/usr/share/doc/apt-xapian-index
/usr/share/doc/apt-xapian-index/changelog.gz
/usr/share/doc/apt-xapian-index/copyright
/usr/share/doc/apt-xapian-index/examples
/usr/share/doc/apt-xapian-index/examples/axi-searchcloud.py.gz
/usr/share/doc/apt-xapian-index/examples/axi-searchasyoutype.py.gz
/usr/share/doc/apt-xapian-index/examples/axi-query-expand.py
/usr/share/doc/apt-xapian-index/examples/axi-query-similar.py
/usr/share/doc/apt-xapian-index/examples/axi-query-adaptivecutoff.py
/usr/share/doc/apt-xapian-index/examples/axi-query-tags.py
/usr/share/doc/apt-xapian-index/examples/axi-query-pkgtype.py.gz
/usr/share/doc/apt-xapian-index/examples/axi-query.py
/usr/share/doc/apt-xapian-index/examples/axi-query-simple.py.gz
/usr/share/doc/apt-xapian-index/examples/aptxapianindex.py.gz
/usr/share/doc/apt-xapian-index/examples/README
/usr/share/doc/apt-xapian-index/README
/usr/share/apt-xapian-index
/usr/share/apt-xapian-index/plugins
/usr/share/apt-xapian-index/plugins/debtags.py
/usr/share/apt-xapian-index/plugins/template.py
/usr/share/apt-xapian-index/plugins/apttags.py
/usr/share/apt-xapian-index/plugins/sizes.py
/usr/share/apt-xapian-index/plugins/descriptions.py
/usr/share/apt-xapian-index/plugins/sections.py
/usr/share/apt-xapian-index/plugins/translated-desc.py
/usr/sbin
/usr/sbin/update-apt-xapian-index
