# var/lib/dpkg/info/araneo-config.list
/.
/etc
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/71-persistent-usb.rules
/etc/udev/rules.d/40-fsl-permissions.rules
/etc/init.d
/etc/init.d/usplash-start
/etc/evtouch
/etc/evtouch/config
/etc/sdio
/etc/sdio/ks7010sd.nss.cfg
/etc/sdio/ks7010sd.cfg
/etc/sdio/ks7010sd.rom
/etc/default
/etc/default/locale
/etc/hal
/etc/hal/fdi
/etc/hal/fdi/policy
/etc/hal/fdi/policy/sd_card.fdi
/etc/event.d
/etc/event.d/setdate
/etc/event.d/ks79xx-sdio
/usr
/usr/share
/usr/share/applications
/usr/share/applications/sharp-electronic-manual.desktop
/usr/share/applications/sharp-user-registration.desktop
/usr/share/usb-modeswitch
/usr/share/usb-modeswitch/C01LC.conf
/usr/share/usb-modeswitch/L-05A.conf
/usr/share/usb-modeswitch/L-05A-script
/usr/share/usb-modeswitch/D12LC-script
/usr/share/usb-modeswitch/L-02A-script
/usr/share/usb-modeswitch/D21LC-script
/usr/share/usb-modeswitch/D11LC.conf
/usr/share/usb-modeswitch/D21LC.conf
/usr/share/usb-modeswitch/C01LC-script
/usr/share/usb-modeswitch/HX003ZT.conf
/usr/share/usb-modeswitch/D31HW.conf
/usr/share/usb-modeswitch/L-05A-2.conf
/usr/share/usb-modeswitch/D12LC.conf
/usr/share/usb-modeswitch/D31HW-script
/usr/share/usb-modeswitch/bmobileZTE-script
/usr/share/usb-modeswitch/bmobileZTE.conf
/usr/share/usb-modeswitch/D11LC-script
/usr/share/usb-modeswitch/HX003ZT-script
/usr/share/usb-modeswitch/L-05A-2-script
/usr/share/usb-modeswitch/L-02A.conf
/usr/share/ume-config-base
/usr/share/ume-config-base/adduser
/usr/share/ume-config-base/adduser/add-user-to-av-groups
/usr/share/gconf
/usr/share/gconf/defaults
/usr/share/gconf/defaults/32_araneo-config
/usr/share/sharp
/usr/share/sharp/MitUser
/usr/share/sharp/MitUser/img
/usr/share/sharp/MitUser/img/Thumbs.db
/usr/share/sharp/MitUser/img/t1_r2.gif
/usr/share/sharp/MitUser/img/t1_rb.gif
/usr/share/sharp/MitUser/img/t1_m3.gif
/usr/share/sharp/MitUser/img/trk_r1_c2.gif
/usr/share/sharp/MitUser/img/trk_r3_c2.gif
/usr/share/sharp/MitUser/img/spacer.gif
/usr/share/sharp/MitUser/img/trk_r2_c2.gif
/usr/share/sharp/MitUser/img/t1_web.gif
/usr/share/sharp/MitUser/img/t1_m1.gif
/usr/share/sharp/MitUser/img/trk_r11_c1.gif
/usr/share/sharp/MitUser/img/t1_m2.gif
/usr/share/sharp/MitUser/img/copyr.gif
/usr/share/sharp/MitUser/img/t1_l2.gif
/usr/share/sharp/MitUser/img/trk_r5_c2.gif
/usr/share/sharp/MitUser/img/btn_rst1.gif
/usr/share/sharp/MitUser/img/t1_lb.gif
/usr/share/sharp/MitUser/img/t1_post.gif
/usr/share/sharp/MitUser/user.ico
/usr/share/sharp/MitUser/user.htm
/usr/share/sharp/MitUser/regist.css
/usr/share/doc
/usr/share/doc/araneo-config
/usr/share/doc/araneo-config/copyright
/usr/share/doc/araneo-config/README.Debian
/usr/share/doc/araneo-config/changelog.gz
/usr/share/oem-config
/usr/share/oem-config/preseed.d
/usr/share/oem-config/preseed.d/locale.seed
/usr/share/oem-config/preseed.d/timezone.seed
/usr/sbin
/usr/bin
