# var/lib/dpkg/info/araneo-config.conffiles
/etc/udev/rules.d/71-persistent-usb.rules
/etc/udev/rules.d/40-fsl-permissions.rules
/etc/init.d/usplash-start
/etc/evtouch/config
/etc/sdio/ks7010sd.nss.cfg
/etc/sdio/ks7010sd.cfg
/etc/sdio/ks7010sd.rom
/etc/default/locale
/etc/hal/fdi/policy/sd_card.fdi
/etc/event.d/setdate
/etc/event.d/ks79xx-sdio
