# var/lib/dpkg/info/at.prerm
#! /bin/sh

if [ -x /etc/init.d/atd ]
then
	invoke-rc.d atd stop || true
fi
