# var/lib/dpkg/info/avahi-daemon.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/avahi-dbus.conf
/etc/default
/etc/default/avahi-daemon
/etc/network
/etc/network/if-up.d
/etc/network/if-up.d/avahi-daemon
/etc/network/if-post-down.d
/etc/init.d
/etc/init.d/avahi-daemon
/etc/avahi
/etc/avahi/avahi-daemon.conf
/etc/avahi/services
/etc/avahi/hosts
/etc/resolvconf
/etc/resolvconf/update-libc.d
/etc/resolvconf/update-libc.d/avahi-daemon
/usr
/usr/share
/usr/share/doc
/usr/share/doc/avahi-daemon
/usr/share/doc/avahi-daemon/examples
/usr/share/doc/avahi-daemon/examples/ssh.service
/usr/share/doc/avahi-daemon/examples/example.service
/usr/share/doc/avahi-daemon/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/avahi-daemon.conf.5.gz
/usr/share/man/man5/avahi.hosts.5.gz
/usr/share/man/man5/avahi.service.5.gz
/usr/share/man/man8
/usr/share/man/man8/avahi-daemon.8.gz
/usr/share/avahi
/usr/share/avahi/introspection
/usr/share/avahi/introspection/RecordBrowser.introspect
/usr/share/avahi/introspection/EntryGroup.introspect
/usr/share/avahi/introspection/Server.introspect
/usr/share/avahi/introspection/AddressResolver.introspect
/usr/share/avahi/introspection/ServiceTypeBrowser.introspect
/usr/share/avahi/introspection/HostNameResolver.introspect
/usr/share/avahi/introspection/ServiceBrowser.introspect
/usr/share/avahi/introspection/ServiceResolver.introspect
/usr/share/avahi/introspection/DomainBrowser.introspect
/usr/share/avahi/avahi-service.dtd
/usr/sbin
/usr/sbin/avahi-daemon
/usr/lib
/usr/lib/avahi
/usr/lib/avahi/avahi-daemon-check-dns.sh
/etc/network/if-post-down.d/avahi-daemon
/usr/share/doc/avahi-daemon/changelog.Debian.gz
/usr/share/doc/avahi-daemon/README
/usr/share/doc/avahi-daemon/NEWS.gz
