# var/lib/dpkg/info/cron.conffiles
/etc/default/cron
/etc/cron.hourly/.placeholder
/etc/init.d/cron
/etc/pam.d/cron
/etc/crontab
/etc/cron.daily/.placeholder
/etc/cron.daily/standard
/etc/cron.monthly/.placeholder
/etc/cron.monthly/standard
/etc/cron.d/.placeholder
/etc/cron.weekly/.placeholder
