# var/lib/dpkg/info/gdm.conffiles
/etc/gdm/gdm.conf-custom
/etc/gdm/PostLogin/Default.sample
/etc/gdm/gdmprefetchlist
/etc/gdm/locale.conf
/etc/gdm/gdm.conf
/etc/gdm/PreSession/Default
/etc/gdm/modules/AccessKeyMouseEvents
/etc/gdm/modules/factory-AccessDwellMouseEvents
/etc/gdm/modules/factory-AccessKeyMouseEvents
/etc/gdm/modules/AccessDwellMouseEvents
/etc/gdm/PostSession/Default
/etc/gdm/XKeepsCrashing
/etc/gdm/Init/Default
/etc/gdm/Xsession
/etc/init.d/gdm
/etc/pam.d/gdm
/etc/pam.d/gdm-autologin
