# var/lib/dpkg/info/man-db.list
/.
/var
/var/cache
/var/cache/man
/etc
/etc/manpath.config
/etc/cron.daily
/etc/cron.daily/man-db
/etc/cron.weekly
/etc/cron.weekly/man-db
/usr
/usr/share
/usr/share/man-db
/usr/share/man-db/chconfig
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/man-db
/usr/share/doc
/usr/share/doc/man-db
/usr/share/doc/man-db/examples
/usr/share/doc/man-db/examples/manpage.example
/usr/share/doc/man-db/examples/manpath.config.gz
/usr/share/doc/man-db/examples/manpage.example.mdoc
/usr/share/doc/man-db/examples/manpage.example.pod
/usr/share/doc/man-db/examples/manpage.example.sgml.gz
/usr/share/doc/man-db/man-db-manual.dvi.gz
/usr/share/doc/man-db/copyright
/usr/share/doc/man-db/changelog.gz
/usr/share/doc/man-db/README.gz
/usr/share/doc/man-db/changelog.Debian.gz
/usr/share/doc/man-db/man-db-manual.txt.gz
/usr/share/doc/man-db/TODO
/usr/share/doc/man-db/NEWS.gz
/usr/share/doc/man-db/THANKS
/usr/share/doc/man-db/man-db-manual.ps.gz
/usr/share/man
/usr/share/man/ja
/usr/share/man/ja/man1
/usr/share/man/ja/man1/apropos.1.gz
/usr/share/man/ja/man1/man.1.gz
/usr/share/man/ja/man1/manpath.1.gz
/usr/share/man/ja/man1/zsoelim.1.gz
/usr/share/man/ja/man1/whatis.1.gz
/usr/share/man/ja/man5
/usr/share/man/ja/man5/manpath.5.gz
/usr/share/man/ja/man8
/usr/share/man/ja/man8/mandb.8.gz
/usr/share/man/ja/man8/catman.8.gz
/usr/share/man/ja/man8/accessdb.8.gz
/usr/share/man/man1
/usr/share/man/man1/lexgrog.1.gz
/usr/share/man/man1/apropos.1.gz
/usr/share/man/man1/man.1.gz
/usr/share/man/man1/manpath.1.gz
/usr/share/man/man1/zsoelim.1.gz
/usr/share/man/man1/manconv.1.gz
/usr/share/man/man1/whatis.1.gz
/usr/share/man/es
/usr/share/man/es/man1
/usr/share/man/es/man1/apropos.1.gz
/usr/share/man/es/man1/man.1.gz
/usr/share/man/es/man1/manpath.1.gz
/usr/share/man/es/man1/zsoelim.1.gz
/usr/share/man/es/man1/whatis.1.gz
/usr/share/man/es/man5
/usr/share/man/es/man5/manpath.5.gz
/usr/share/man/es/man8
/usr/share/man/es/man8/mandb.8.gz
/usr/share/man/es/man8/catman.8.gz
/usr/share/man/es/man8/accessdb.8.gz
/usr/share/man/de
/usr/share/man/de/man1
/usr/share/man/de/man1/apropos.1.gz
/usr/share/man/de/man1/man.1.gz
/usr/share/man/de/man1/manpath.1.gz
/usr/share/man/de/man1/zsoelim.1.gz
/usr/share/man/de/man1/whatis.1.gz
/usr/share/man/de/man5
/usr/share/man/de/man5/manpath.5.gz
/usr/share/man/de/man8
/usr/share/man/de/man8/mandb.8.gz
/usr/share/man/de/man8/catman.8.gz
/usr/share/man/de/man8/accessdb.8.gz
/usr/share/man/it
/usr/share/man/it/man1
/usr/share/man/it/man1/apropos.1.gz
/usr/share/man/it/man1/man.1.gz
/usr/share/man/it/man1/manpath.1.gz
/usr/share/man/it/man1/zsoelim.1.gz
/usr/share/man/it/man1/whatis.1.gz
/usr/share/man/it/man5
/usr/share/man/it/man5/manpath.5.gz
/usr/share/man/it/man8
/usr/share/man/it/man8/mandb.8.gz
/usr/share/man/it/man8/catman.8.gz
/usr/share/man/it/man8/accessdb.8.gz
/usr/share/man/fr
/usr/share/man/fr/man1
/usr/share/man/fr/man1/lexgrog.1.gz
/usr/share/man/fr/man1/apropos.1.gz
/usr/share/man/fr/man1/man.1.gz
/usr/share/man/fr/man1/manpath.1.gz
/usr/share/man/fr/man1/zsoelim.1.gz
/usr/share/man/fr/man1/whatis.1.gz
/usr/share/man/fr/man5
/usr/share/man/fr/man5/manpath.5.gz
/usr/share/man/fr/man8
/usr/share/man/fr/man8/mandb.8.gz
/usr/share/man/fr/man8/catman.8.gz
/usr/share/man/fr/man8/accessdb.8.gz
/usr/share/man/man5
/usr/share/man/man5/manpath.5.gz
/usr/share/man/man8
/usr/share/man/man8/mandb.8.gz
/usr/share/man/man8/catman.8.gz
/usr/share/man/man8/accessdb.8.gz
/usr/share/doc-base
/usr/share/doc-base/man-db
/usr/share/bug
/usr/share/bug/man-db
/usr/share/bug/man-db/presubj
/usr/sbin
/usr/sbin/accessdb
/usr/bin
/usr/bin/whatis
/usr/bin/manpath
/usr/bin/catman
/usr/bin/zsoelim
/usr/bin/lexgrog
/usr/lib
/usr/lib/man-db
/usr/lib/man-db/globbing
/usr/lib/man-db/man
/usr/lib/man-db/mandb
/usr/lib/man-db/manconv
/usr/bin/apropos
