# var/lib/dpkg/info/network-manager.conffiles
/etc/init.d/NetworkManager
/etc/NetworkManager/nm-system-settings.conf
/etc/NetworkManager/dispatcher.d/01ifupdown
/etc/dbus-1/system.d/nm-system-settings.conf
/etc/dbus-1/system.d/nm-dhcp-client.conf
/etc/dbus-1/system.d/nm-avahi-autoipd.conf
/etc/dbus-1/system.d/NetworkManager.conf
/etc/dbus-1/system.d/nm-dispatcher.conf
