# var/lib/dpkg/info/network-manager.prerm
#!/bin/sh

set -e

case "$1" in
  	remove)
		RESULT=0
                invoke-rc.d NetworkManager stop
	;;
esac



exit 0
