# var/lib/dpkg/info/oem-config.list
/.
/etc
/etc/init.d
/etc/init.d/oem-config
/var
/var/lib
/var/lib/localechooser
/usr
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_oem-config.py
/usr/share/pyshared-data
/usr/share/pyshared-data/oem-config
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/oem-config
/usr/share/doc
/usr/share/doc/oem-config
/usr/share/doc/oem-config/copyright
/usr/share/doc/oem-config/changelog.gz
/usr/share/oem-config
/usr/share/oem-config/pixmaps
/usr/share/oem-config/pixmaps/timezone
/usr/share/oem-config/pixmaps/timezone/timezone_1.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_-5.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_13.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_12.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_0.0.png
/usr/share/oem-config/pixmaps/timezone/cc.png
/usr/share/oem-config/pixmaps/timezone/timezone_5.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_9.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_6.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_4.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_10.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_-9.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_-3.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_7.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_3.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_6.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_10.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_-10.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_11.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_2.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_-4.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_9.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_-2.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_11.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_5.75.png
/usr/share/oem-config/pixmaps/timezone/timezone_3.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_-9.0.png
/usr/share/oem-config/pixmaps/timezone/bg.png
/usr/share/oem-config/pixmaps/timezone/timezone_-6.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_8.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_5.5.png
/usr/share/oem-config/pixmaps/timezone/timezone_-1.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_-3.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_-7.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_12.75.png
/usr/share/oem-config/pixmaps/timezone/timezone_-8.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_4.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_-11.0.png
/usr/share/oem-config/pixmaps/timezone/timezone_-5.0.png
/usr/share/oem-config/pixmaps/world_map-960.png
/usr/share/oem-config/preseed.d
/usr/lib
/usr/lib/oem-config
/usr/lib/oem-config/timezone
/usr/lib/oem-config/timezone/tzsetup
/usr/lib/oem-config/timezone/post-base-installer
/usr/lib/oem-config/timezone/tzsetup-wrapper
/usr/lib/oem-config/timezone/tzmap
/usr/lib/oem-config/console
/usr/lib/oem-config/console/console-setup.postinst
/usr/lib/oem-config/console/console-setup-apply
/usr/lib/oem-config/user
/usr/lib/oem-config/user/user-setup
/usr/lib/oem-config/user/user-setup-ask
/usr/lib/oem-config/user/functions.sh
/usr/lib/oem-config/user/user-setup-wrapper
/usr/lib/oem-config/user/user-setup-apply
/usr/lib/oem-config/user/reserved-usernames
/usr/lib/oem-config/compat
/usr/lib/oem-config/compat/udpkg
/usr/lib/oem-config/compat/apt-install
/usr/lib/oem-config/compat/anna-install
/usr/lib/oem-config/welcome
/usr/lib/oem-config/welcome/welcome
/usr/lib/oem-config/network
/usr/lib/oem-config/network/netcfg
/usr/lib/oem-config/network/netcfg-wrapper
/usr/lib/oem-config/oem_config
/usr/lib/oem-config/oem_config/keyboard_names.py
/usr/lib/oem-config/oem_config/components
/usr/lib/oem-config/oem_config/components/network.py
/usr/lib/oem-config/oem_config/components/timezone.py
/usr/lib/oem-config/oem_config/components/welcome.py
/usr/lib/oem-config/oem_config/components/console_setup_apply.py
/usr/lib/oem-config/oem_config/components/language.py
/usr/lib/oem-config/oem_config/components/language_apply.py
/usr/lib/oem-config/oem_config/components/tasks.py
/usr/lib/oem-config/oem_config/components/console_setup.py
/usr/lib/oem-config/oem_config/components/__init__.py
/usr/lib/oem-config/oem_config/components/user.py
/usr/lib/oem-config/oem_config/components/timezone_apply.py
/usr/lib/oem-config/oem_config/timezone_map.py
/usr/lib/oem-config/oem_config/debconffilter.py
/usr/lib/oem-config/oem_config/osextras.py
/usr/lib/oem-config/oem_config/i18n.py
/usr/lib/oem-config/oem_config/__init__.py
/usr/lib/oem-config/oem_config/frontend
/usr/lib/oem-config/oem_config/frontend/__init__.py
/usr/lib/oem-config/oem_config/frontend/base.py
/usr/lib/oem-config/oem_config/cairo_extensions.py
/usr/lib/oem-config/oem_config/im_switch.py
/usr/lib/oem-config/oem_config/tz.py
/usr/lib/oem-config/oem_config/filteredcommand.py
/usr/lib/oem-config/language
/usr/lib/oem-config/language/localechooser-wrapper
/usr/lib/oem-config/language/translation-check
/usr/lib/oem-config/language/regionmap
/usr/lib/oem-config/language/shortlists
/usr/lib/oem-config/language/localechooser
/usr/lib/oem-config/language/languagelist
/usr/lib/oem-config/language/languagemap
/usr/lib/oem-config/language/localechooser-apply
/usr/lib/oem-config/language/languagelist.data.gz
/usr/sbin
/usr/sbin/oem-config-firstboot
/usr/sbin/oem-config-prepare
/usr/sbin/oem-config-dm
/usr/sbin/oem-config
