# var/lib/dpkg/info/ppp.conffiles
/etc/init.d/pppd-dns
/etc/chatscripts/pap
/etc/bash_completion.d/pon
/etc/pam.d/ppp
/etc/apm/event.d/ppp
/etc/logrotate.d/ppp
/etc/ppp/ip-up
/etc/ppp/ip-up.d/0000usepeerdns
/etc/ppp/pppoe_on_boot
/etc/ppp/ipv6-up
/etc/ppp/ip-down.d/0000usepeerdns
/etc/ppp/ipv6-down
/etc/ppp/options
/etc/ppp/ip-down
