# var/lib/dpkg/info/pulseaudio.conffiles
/etc/pulse/default.pa
/etc/pulse/system.pa
/etc/pulse/daemon.conf
/etc/X11/Xsession.d/70pulseaudio
/etc/init.d/pulseaudio
/etc/default/pulseaudio
