# var/lib/dpkg/info/readahead.conffiles
/etc/readahead/boot
/etc/readahead/desktop
/etc/init.d/stop-readahead
/etc/init.d/readahead-desktop
/etc/init.d/readahead
