# var/lib/dpkg/info/readahead.list
/.
/sbin
/sbin/readahead-watch
/sbin/readahead-list
/etc
/etc/readahead
/etc/readahead/boot
/etc/readahead/desktop
/etc/init.d
/etc/init.d/stop-readahead
/etc/init.d/readahead-desktop
/etc/init.d/readahead
/usr
/usr/share
/usr/share/doc
/usr/share/doc/readahead
/usr/share/doc/readahead/copyright
/usr/share/doc/readahead/changelog.gz
/usr/share/doc/readahead/changelog.Debian.gz
/usr/share/doc/readahead/README.Debian
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/readahead-list.8.gz
/usr/share/man/man8/readahead-watch.8.gz
