# var/lib/dpkg/info/system-tools-backends.list
/.
/etc
/etc/init.d
/etc/init.d/system-tools-backends
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/system-tools-backends.conf
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/system-tools-backends.1.gz
/usr/share/system-tools-backends-2.0
/usr/share/system-tools-backends-2.0/scripts
/usr/share/system-tools-backends-2.0/scripts/Time
/usr/share/system-tools-backends-2.0/scripts/Time/TimeDate.pm
/usr/share/system-tools-backends-2.0/scripts/Time/NTP.pm
/usr/share/system-tools-backends-2.0/scripts/IfacesConfig.pm
/usr/share/system-tools-backends-2.0/scripts/NTPConfig.pm
/usr/share/system-tools-backends-2.0/scripts/NFSConfig.pm
/usr/share/system-tools-backends-2.0/scripts/SMBConfig.pm
/usr/share/system-tools-backends-2.0/scripts/Platform.pm
/usr/share/system-tools-backends-2.0/scripts/UserConfig.pm
/usr/share/system-tools-backends-2.0/scripts/GroupsConfig.pm
/usr/share/system-tools-backends-2.0/scripts/Shares
/usr/share/system-tools-backends-2.0/scripts/Shares/SMB.pm
/usr/share/system-tools-backends-2.0/scripts/Shares/NFS.pm
/usr/share/system-tools-backends-2.0/scripts/HostsConfig.pm
/usr/share/system-tools-backends-2.0/scripts/Utils
/usr/share/system-tools-backends-2.0/scripts/Utils/Monitor.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/File.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/Backend.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/Platform.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/Replace.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/XML.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/Report.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/Util.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/DBus.pm
/usr/share/system-tools-backends-2.0/scripts/Utils/Parse.pm
/usr/share/system-tools-backends-2.0/scripts/Users
/usr/share/system-tools-backends-2.0/scripts/Users/Users.pm
/usr/share/system-tools-backends-2.0/scripts/Users/Shells.pm
/usr/share/system-tools-backends-2.0/scripts/Users/Groups.pm
/usr/share/system-tools-backends-2.0/scripts/StbObject.pm
/usr/share/system-tools-backends-2.0/scripts/Network
/usr/share/system-tools-backends-2.0/scripts/Network/Ifaces.pm
/usr/share/system-tools-backends-2.0/scripts/Network/Hosts.pm
/usr/share/system-tools-backends-2.0/scripts/UsersConfig.pm
/usr/share/system-tools-backends-2.0/scripts/Init
/usr/share/system-tools-backends-2.0/scripts/Init/Services.pm
/usr/share/system-tools-backends-2.0/scripts/Init/ServicesList.pm
/usr/share/system-tools-backends-2.0/scripts/SystemToolsBackends.pl
/usr/share/system-tools-backends-2.0/scripts/TimeConfig.pm
/usr/share/system-tools-backends-2.0/scripts/ServicesConfig.pm
/usr/share/system-tools-backends-2.0/files
/usr/share/system-tools-backends-2.0/files/general_isdn_ppp_options
/usr/share/system-tools-backends-2.0/files/general_pppoe_ppp_options
/usr/share/system-tools-backends-2.0/files/general_gprs_chatscript
/usr/share/PolicyKit
/usr/share/PolicyKit/policy
/usr/share/PolicyKit/policy/system-tools-backends.policy
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.GroupsConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.Platform.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.UsersConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.UserConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.ServicesConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.SMBConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.NFSConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.TimeConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.HostsConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.IfacesConfig.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.service
/usr/share/dbus-1/system-services/org.freedesktop.SystemToolsBackends.NTPConfig.service
/usr/share/doc
/usr/share/doc/system-tools-backends
/usr/share/doc/system-tools-backends/copyright
/usr/share/doc/system-tools-backends/README
/usr/share/doc/system-tools-backends/AUTHORS
/usr/share/doc/system-tools-backends/NEWS.gz
/usr/share/doc/system-tools-backends/changelog.Debian.gz
/usr/bin
/usr/bin/system-tools-backends
