# var/lib/dpkg/info/udev.conffiles
/etc/scsi_id.config
/etc/init.d/udev-finish
/etc/init.d/udev
/etc/udev/udev.conf
/etc/udev/rules.d/README
