# var/lib/dpkg/info/ufw.conffiles
/etc/init.d/ufw
/etc/ufw/sysctl.conf
/etc/default/ufw
