# var/lib/dpkg/info/ufw.list
/.
/etc
/etc/init.d
/etc/init.d/ufw
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/sysctl.conf
/etc/default
/etc/default/ufw
/usr
/usr/share
/usr/share/doc
/usr/share/doc/ufw
/usr/share/doc/ufw/README.gz
/usr/share/doc/ufw/README.Debian
/usr/share/doc/ufw/changelog.gz
/usr/share/doc/ufw/changelog.pre-0.25.gz
/usr/share/doc/ufw/TODO.Debian
/usr/share/doc/ufw/copyright
/usr/share/doc/ufw/examples
/usr/share/doc/ufw/examples/initscript.example
/usr/share/doc/ufw/examples/skel-ui.example
/usr/share/doc/ufw/changelog.Debian.gz
/usr/share/ufw
/usr/share/ufw/ufw-init
/usr/share/ufw/ufw-init-functions
/usr/share/ufw/after.rules
/usr/share/ufw/after6.rules
/usr/share/ufw/before.rules
/usr/share/ufw/before6.rules
/usr/share/ufw/ufw.conf
/usr/share/ufw/user.rules
/usr/share/ufw/user6.rules
/usr/share/ufw/after.rules.md5sum
/usr/share/ufw/after6.rules.md5sum
/usr/share/ufw/before.rules.md5sum
/usr/share/ufw/before6.rules.md5sum
/usr/share/ufw/user.rules.md5sum
/usr/share/ufw/user6.rules.md5sum
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ufw.8.gz
/usr/share/pyshared
/usr/share/pyshared/ufw
/usr/share/pyshared/ufw/__init__.py
/usr/share/pyshared/ufw/backend.py
/usr/share/pyshared/ufw/backend_iptables.py
/usr/share/pyshared/ufw/common.py
/usr/share/pyshared/ufw/frontend.py
/usr/share/pyshared/ufw/util.py
/usr/share/pyshared/ufw/applications.py
/usr/share/pyshared/ufw-0.27_0ubuntu2.egg-info
/usr/share/pyshared-data
/usr/share/pyshared-data/ufw
/usr/lib
/usr/lib/python2.6
/usr/lib/python2.6/dist-packages
/usr/lib/python2.6/dist-packages/ufw
/usr/lib/python2.5
/usr/lib/python2.5/site-packages
/usr/lib/python2.5/site-packages/ufw
/usr/sbin
/usr/sbin/ufw
/var
/var/lib
/var/lib/ufw
/var/lib/ufw/messages
/var/lib/ufw/messages/ufw.pot
/usr/lib/python2.6/dist-packages/ufw-0.27_0ubuntu2.egg-info
/usr/lib/python2.6/dist-packages/ufw/__init__.py
/usr/lib/python2.6/dist-packages/ufw/backend.py
/usr/lib/python2.6/dist-packages/ufw/backend_iptables.py
/usr/lib/python2.6/dist-packages/ufw/common.py
/usr/lib/python2.6/dist-packages/ufw/frontend.py
/usr/lib/python2.6/dist-packages/ufw/util.py
/usr/lib/python2.6/dist-packages/ufw/applications.py
/usr/lib/python2.5/site-packages/ufw-0.27_0ubuntu2.egg-info
/usr/lib/python2.5/site-packages/ufw/__init__.py
/usr/lib/python2.5/site-packages/ufw/backend.py
/usr/lib/python2.5/site-packages/ufw/backend_iptables.py
/usr/lib/python2.5/site-packages/ufw/common.py
/usr/lib/python2.5/site-packages/ufw/frontend.py
/usr/lib/python2.5/site-packages/ufw/util.py
/usr/lib/python2.5/site-packages/ufw/applications.py
