# var/lib/dpkg/info/update-motd.list
/.
/etc
/etc/cron.d
/etc/cron.d/update-motd
/etc/update-motd.d
/etc/update-motd.d/hourly
/etc/update-motd.d/daily
/etc/update-motd.d/weekly
/etc/update-motd.d/monthly
/var
/var/lib
/var/lib/update-motd
/usr
/usr/share
/usr/share/doc
/usr/share/doc/update-motd
/usr/share/doc/update-motd/changelog.gz
/usr/share/doc/update-motd/README
/usr/share/doc/update-motd/THANKS
/usr/share/doc/update-motd/copyright
/usr/share/doc/update-motd/changelog.Debian.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/update-motd.1.gz
/usr/sbin
/usr/sbin/update-motd
