# var/lib/dpkg/info/usplash.list
/.
/sbin
/sbin/usplash
/sbin/usplash_down
/sbin/usplash_write
/lib
/lib/init
/lib/init/splash-functions
/etc
/etc/kbd
/etc/kbd/config.d
/etc/kbd/config.d/usplash
/etc/console-tools
/etc/console-tools/config.d
/etc/console-tools/config.d/usplash
/etc/init.d
/etc/init.d/usplash
/usr
/usr/share
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-top
/usr/share/initramfs-tools/scripts/init-top/usplash
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/usplash
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/usplash.py
/usr/share/doc
/usr/share/doc/usplash
/usr/share/doc/usplash/copyright
/usr/share/doc/usplash/changelog.gz
/usr/share/doc/usplash/README
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/usplash.8.gz
/usr/share/man/man8/usplash_write.8.gz
/usr/sbin
/usr/sbin/update-usplash-theme
