# var/lib/dpkg/info/usplash.conffiles
/etc/kbd/config.d/usplash
/etc/console-tools/config.d/usplash
/etc/init.d/usplash
