# var/lib/dpkg/info/wpasupplicant.conffiles
/etc/dbus-1/system.d/wpa_supplicant.conf
/etc/wpa_supplicant/action_wpa.sh
/etc/wpa_supplicant/ifupdown.sh
/etc/wpa_supplicant/functions.sh
/etc/logrotate.d/wpa_supplicant
/etc/logrotate.d/wpa_action
/etc/init.d/wpa-ifupdown
