# var/lib/dpkg/info/xserver-xorg-input-evtouch.list
/.
/etc
/etc/init.d
/etc/init.d/xserver-xorg-input-evtouch
/etc/evtouch
/usr
/usr/share
/usr/share/xf86-input-evtouch
/usr/share/xf86-input-evtouch/empty_cursor.xbm
/usr/share/applications
/usr/share/applications/calibrate.desktop
/usr/share/pixmaps
/usr/share/pixmaps/calibrate.png
/usr/share/doc
/usr/share/doc/xserver-xorg-input-evtouch
/usr/share/doc/xserver-xorg-input-evtouch/copyright
/usr/share/doc/xserver-xorg-input-evtouch/README
/usr/share/doc/xserver-xorg-input-evtouch/README.Debian
/usr/share/doc/xserver-xorg-input-evtouch/evtouch.html
/usr/share/doc/xserver-xorg-input-evtouch/README.calibration
/usr/share/doc/xserver-xorg-input-evtouch/TODO
/usr/share/doc/xserver-xorg-input-evtouch/changelog.gz
/usr/share/doc/xserver-xorg-input-evtouch/changelog.Debian.gz
/usr/share/hal
/usr/share/hal/fdi
/usr/share/hal/fdi/policy
/usr/share/hal/fdi/policy/20thirdparty
/usr/share/hal/fdi/policy/20thirdparty/50-Sharp.fdi
/usr/share/hal/fdi/policy/20thirdparty/50-itm.fdi
/usr/share/hal/fdi/policy/20thirdparty/50-Fujitsu.fdi
/usr/share/hal/fdi/policy/20thirdparty/50-eGalax.fdi
/usr/share/hal/fdi/policy/20thirdparty/50-touchkit.fdi
/usr/share/hal/fdi/policy/20thirdparty/50-ideaco.fdi
/usr/share/hal/fdi/policy/20thirdparty/50-Panasonic.fdi
/usr/lib
/usr/lib/xf86-input-evtouch
/usr/lib/xf86-input-evtouch/ev_calibrate
/usr/lib/xf86-input-evtouch/calibrate.sh
/usr/lib/xorg
/usr/lib/xorg/modules
/usr/lib/xorg/modules/input
/usr/lib/xorg/modules/input/evtouch_drv.so
/usr/bin
/usr/bin/ev_calibrate
/usr/bin/calibrate_touchscreen
