# var/lib/dpkg/info/xserver-xorg-input-evtouch.postinst
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/xserver-xorg-input-evtouch" ]; then
	update-rc.d xserver-xorg-input-evtouch start 25 2 3 4 5 . >/dev/null || exit $?
fi
# End automatically added section
