# etc/cron.weekly/apt-xapian-index
#!/bin/sh

CMD=/usr/sbin/update-apt-xapian-index

# Rebuild the index
if [ -x $CMD ]
then
	nice ionice -c3	$CMD --quiet
fi
