# etc/xdg/autostart/at-spi-registryd-wrapper.desktop
[Desktop Entry]
Type=Application
Name=AT SPI Registry Wrapper
Exec=/usr/lib/gnome-session/helpers/at-spi-registryd-wrapper
OnlyShowIn=GNOME;
NoDisplay=true
AutostartCondition=GNOME /desktop/gnome/interface/accessibility
X-GNOME-Autostart-Phase=Initialization
X-GNOME-Bugzilla-Bugzilla=GNOME
X-GNOME-Bugzilla-Product=gnome-session
X-GNOME-Bugzilla-Component=general
X-GNOME-Bugzilla-Version=2.26.1
X-Ubuntu-Gettext-Domain=gnome-session-2.0
