# etc/xdg/autostart/fix-kbd-layout.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Set keyboard map	
Exec=setxkbmap
Terminal=false
Type=Application
Categories=
NoDisplay=true
