# etc/xdg/autostart/touchcruiser-indicator.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Touch Cruiser Indicater
Name[ja_JP]=タッチクルーザーインジケータ
Comment=I display the mode of the touch cruiser.
Comment[ja_JP]=タッチクルーザーのモードを表示します
Icon=scrollmode
Exec=/usr/bin/touchcruiser-indicator
Terminal=false
Type=Application
Categories=
OnlyShowIn=GNOME;
X-GNOME-Autostart-enabled=true
