# etc/xdg/autostart/user-first-login.desktop
[Desktop Entry]
Name=First Login
Exec=user-first-login
Terminal=false
NoDisplay=true
Type=Application
