# etc/update-motd.d/90-updates-available
/usr/lib/update-notifier/update-motd-updates-available